# playbooks/files/runit/minidlna/run
#!/usr/bin/env bash

# Ejecuto minidlna
exec chpst -u cubie:cubie minidlnad -d -f /etc/minidlna.conf
